# dist/rpm/mariner-usb-gadget.service
[Unit]
Description=Mariner USB Mass Storage Gadget
Documentation=https://github.com/amd989/mariner
DefaultDependencies=no
After=local-fs.target
Before=mariner3d.service
ConditionPathExists=/piusb.bin

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/mariner3d/gadget-start
ExecStop=/usr/lib/mariner3d/gadget-stop

[Install]
WantedBy=multi-user.target
